# support/cloud-nginx-proxy.service.erb
[Unit]
Description=nginx reverse proxy reconfiguration
After=network-online.target etcd2.service

[Service]
SyslogIdentifier=cloud-proxy-configuration
User=www-data
ExecStartPre=/opt/mee/cloud-nginx-proxy/bin/catch-up --lb-dir <%= @nginx.dir %> --etcd-prefix <%= @etcd.watch_base %> --notify-key <%= @nginx.notify_key%>
ExecStart=/opt/coreos/bin/etcdctl exec-watch <%= @etcd.watch_base %> --recursive -- /opt/mee/cloud-nginx-proxy/bin/watcher --lb-dir <%= @nginx.dir %> --etcd-prefix <%= @etcd.watch_base %> --notify-key <%= @nginx.notify_key%>
Nice=15
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
